# photo_frame_server.service
[Unit]
Description=Digital Photo Frame Viewer Server
After=network.target photo_frame_sync.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/digital_photo_frame
Environment="PORT=8000"
ExecStart=/home/pi/digital_photo_frame/venv/bin/python /home/pi/digital_photo_frame/viewer_server.py /home/pi/digital_photo_frame/config_synology.yaml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
